# programs_VF61_VF50/Meter.service
[Unit]
Description=Smart Meter
After=multi-user.target
 
[Service]
Type=simple
ExecStart=/home/root/init.sh
 
[Install]
WantedBy=multi-user.target